# devops/templates/supervisor.ini
[program:{{ app_name }}]
command=/usr/local/bin/uwsgi
  --socket /tmp/{{ app_name }}.sock
  --logto {{ webapps_dir }}/{{ app_name }}/log/application.log
  --home {{ webapps_dir }}/{{ app_name }}/venv
  --pythonpath {{ webapps_dir }}/{{ app_name }}/src
  --wsgi-file {{ webapps_dir }}/{{ app_name }}/src/{{ wsgi_file }}
  --callable {{ wsgi_callable }}
  --max-requests 1000
  --master
  --processes 1
  --chmod
directory={{ webapps_dir }}/{{ app_name }}/src
autostart=true
autorestart=true

[program:{{ app_name }}-celery]
command={{ webapps_dir }}/{{ app_name }}/venv/bin/celery
  worker
  --app={{ app_name }}.tasks
  --loglevel=INFO
directory={{ webapps_dir }}/{{ app_name }}/src
user=nobody
numprocs=1
stdout_logfile={{ webapps_dir }}/{{ app_name }}/log/celery-worker.log
stderr_logfile={{ webapps_dir }}/{{ app_name }}/log/celery-worker.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
killasgroup=true
priority=1000